# crontab/crontab_cactus
SHELL=/bin/bash -l

#00 09 * * * qsub -o /lfs/h2/emc/ptmp/steven.simon/evs_out/jevs_global_ens_wave_grid2obs_prep.out /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/EVS/dev/drivers/scripts/global_ens/prep/jevs_global_ens_wave_grid2obs_prep.sh
#00 10 * * * qsub -o /lfs/h2/emc/ptmp/steven.simon/evs_out/jevs_global_ens_wave_grid2obs_stats.out /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/EVS/dev/drivers/scripts/global_ens/stats/jevs_global_ens_wave_grid2obs_stats.sh
#00 11 * * * qsub -o /lfs/h2/emc/ptmp/steven.simon/evs_out/jevs_global_ens_wave_grid2obs_plots.out /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/EVS/dev/drivers/scripts/global_ens/plots/jevs_global_ens_wave_grid2obs_plots.sh

00 09 * * * sh /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/SMSIMON/trigger_scripts/global_ens_wave_prep_trigger_script.sh
00 10 * * * sh /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/SMSIMON/trigger_scripts/global_ens_wave_stats_trigger_script.sh
00 11 * * * sh /lfs/h2/emc/vpppg/noscrub/steven.simon/gitworkspace/SMSIMON/trigger_scripts/global_ens_wave_plots_trigger_script.sh


